# private_dot_config/systemd/user/wallpaper-apply.timer
[Unit]
Description=Run wallpaper apply every 5 minutes

[Timer]
OnStartupSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
