# private_dot_config/systemd/user/wallpaper-update.timer
[Unit]
Description=Run wallpaper update every (hour % 6) == 0

[Timer]
OnCalendar=00/6:00:00
Persistent=True

[Install]
WantedBy=timers.target
