# private_dot_config/systemd/user/wallpaper-update.service
[Unit]
Description=Run wallpaper update
After=network-online.target nss-lookup.target

[Service]
ExecStart=%h/.local/bin/wallpaper update
