# private_dot_config/systemd/user/wallpaper-update.timer
[Unit]
Description=Run wallpaper update every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
