# private_dot_config/systemd/user/wallpaper-apply.service
[Unit]
Description=Run wallpaper apply

[Service]
ExecStart=%h/.local/bin/wallpaper apply
