# private_dot_config/systemd/user/wallpaper-update.service
[Unit]
Description=Run wallpaper update

[Service]
ExecStart=%h/.local/bin/wallpaper update
